# launchd/templates/ai.openviking.mcp.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openviking.mcp</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/projects/openviking-local/.venv/bin/python3</string>
        <string>__HOME__/projects/openviking-local/server.py</string>
        <string>--port</string>
        <string>__OPENVIKING_PORT__</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>__HOME__/.openclaw/openviking-data</string>
    <key>StandardOutPath</key>
    <string>__HOME__/.openclaw/logs/openviking.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.openclaw/logs/openviking.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>__HOMEBREW_PREFIX__/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
