# launchd/templates/ai.openclaw.gateway.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.gateway</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-i</string>
        <string>__HOMEBREW_PREFIX__/bin/openclaw</string>
        <string>gateway</string>
        <string>--port</string>
        <string>__GATEWAY_PORT__</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>__HOME__/.openclaw/logs/gateway.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.openclaw/logs/gateway.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>__HOMEBREW_PREFIX__/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>NODE_OPTIONS</key>
        <string>--max-old-space-size=512</string>
        <key>HTTPS_PROXY</key>
        <string>__HTTP_PROXY__</string>
        <key>HTTP_PROXY</key>
        <string>__HTTP_PROXY__</string>
    </dict>
</dict>
</plist>
